# etc/cron/document_system.cron
#!/bin/bash

cd /home/pschaaf

OUTPUT=local/sync/system-docs/documentation_`hostname`.txt
TMPFILE=`mktemp /tmp/$$.XXXXXX` || exit 1
rm -f $OUTPUT
chown pschaaf $TMPFILE
bin/document_system -q 2>&1 >> $TMPFILE
mv $TMPFILE $OUTPUT
